# build/supervisord.conf
[supervisord]
nodaemon=true

[inet_http_server]
port=127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[program:monitor]
command=node index.js
autostart=true
directory=/usr/src/monitor
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:relayer]
command=node -r /usr/src/monitor/env.js app.js
autostart=false
directory=/usr/src/app/relayer
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:redis]
command=redis-server
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/var/log/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/%(program_name)s_stdout.log
stdout_logfile_maxbytes=10MB
